# templates/etc/profile.d/conman.sh.erb
######### Managed by Puppet #########

export CONMAN_HOST=<%= @conman_server %>
export CONMAN_PORT=<%= @conman_port %>

