# templates/etc/profile.d/conman.csh.erb
######### Managed by Puppet #########
# <%= scope['settings::server'] %>:<%= @__file__ %>
setenv CONMAN_HOST <%= @conman_server %>
setenv CONMAN_PORT <%= @conman_port %>

